# vc_software/systemd/vc-ota.service
# vc-ota.service 파일

[Unit]
Description=Vehicle Computer OTA Update Service
After=network.target

[Service]
WorkingDirectory=/home/hj/vc_software/ota_updater
ExecStart=/usr/bin/python3 /home/hj/vc_software/ota_updater/ota_service.py

Type=simple
User=hj
Group=hj
Restart=always

[Install]
WantedBy=multi-user.target